# templates/hipache_upstart.conf
description     "Run hipache"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

respawn

script
    hipache --config /etc/hipache/config.json
end script